# deploy/newapi-gateway.service
[Unit]
Description=NewAPI Gateway Service
Documentation=https://github.com/Calcium-Ion/new-api
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/newapi-gateway

# Executable path
ExecStart=/opt/newapi-gateway/gateway-aggregator --port 3030 --log-dir /opt/newapi-gateway/logs

# Restart policy
Restart=always
RestartSec=5s
StartLimitInterval=0

# Environment variables
Environment="GIN_MODE=release"
Environment="PORT=3030"
Environment="LOG_DIR=/opt/newapi-gateway/logs"
Environment="SESSION_SECRET=CHANGE_THIS_TO_RANDOM_SECRET"

# Optional: database configuration (if needed)
# Environment="SQL_DSN=sqlite:///opt/newapi-gateway/data/newapi.db"
# Environment="SQL_DSN=user:password@tcp(localhost:3306)/newapi?charset=utf8mb4&parseTime=True&loc=Local"

# Optional: Redis configuration
# Environment="REDIS_CONN_STRING=redis://localhost:6379"
# Environment="REDIS_PASSWORD=your_redis_password"

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
MemoryLimit=1G
CPUQuota=100%

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=newapi-gateway

[Install]
WantedBy=multi-user.target
